# package/system/librerouter-hw/files/etc/uci-defaults/99-run-librerouter-hw-check
#!/bin/sh

# we wait until the LR is configured, up and stable
(sleep 200 && librerouter-hw-check )&

exit 0
